# services.d/server/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting server"
exec node ./server/dist/server/index.js