# services.d/client/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting client"
exec ./node_modules/.bin/vite