# services.d/client/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting client"
exec /usr/app/node_modules/.bin/vite